# aria-ceo.service
[Unit]
Description=Aria CEO v7.0 Config-Driven Multi-Agent Orchestrator
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/aria-system
ExecStart=/usr/bin/python3 /opt/aria-system/agents/aria_ceo_v7.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
